# deploy/rmpg-webhook.service
[Unit]
Description=RMPG Flex GitHub Webhook Listener
Documentation=https://github.com/your-repo/rmpg-flex
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node /opt/rmpg-flex/deploy/webhook-listener.js
WorkingDirectory=/opt/rmpg-flex
Restart=always
RestartSec=5
Environment=WEBHOOK_PORT=9000
Environment=REPO_DIR=/opt/rmpg-flex
Environment=DEPLOY_SCRIPT=/opt/deploy-rmpg.sh

# Security hardening
NoNewPrivileges=false
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rmpg-webhook

[Install]
WantedBy=multi-user.target
